# bin/empyr-a01.service
[Unit]
Description=Empyr A01 Web Service
StartLimitIntervalSec=0
After=network-online.target

[Service]
Type=simple
User=empyr
Group=empyr
PIDFile=/run/empyr-a01.pid
WorkingDirectory=/var/www/a01/bin
Environment=EMPYR_ENV=production
ExecStart=/var/www/a01/bin/empyr start server
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure
RestartSec=3600

[Install]
WantedBy=multi-user.target
